# systemd/questix_robot_manager.service
[Unit]
Description=Questix Robot Manager Web UI
After=network.target

[Service]
Type=simple
# Change User/Group to match your system username
User=ubuntu
Group=ubuntu
ExecStart=/usr/bin/python3 -m uvicorn robot_manager.app:app --host 127.0.0.1 --port 8888
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=questix_robot_manager

[Install]
WantedBy=multi-user.target
