# systemd/questix_robot.service
[Unit]
Description=Questix Robot ROS2 Launch Service
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
# Change User/Group to match your system username
User=ubuntu
Group=ubuntu
SupplementaryGroups=gpio input dialout
EnvironmentFile=-/etc/questix_robot/launch.env
# Defaults — override via launch.env or questix_robot.env
Environment=ROS_DISTRO=jazzy
Environment=ROBOT_WS=/home/ubuntu/robot_ws
Environment=HOME=/home/ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/opt/questix_robot/questix_robot_launcher.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=questix_robot

# Graceful shutdown
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
