# ansible/roles/robot_autostart/templates/questix_robot.service.j2
[Unit]
Description=Questix Robot ROS2 Launch Service
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User={{ target_user }}
Group={{ target_user }}
SupplementaryGroups=gpio input dialout
EnvironmentFile=-/etc/questix_robot/launch.env
Environment=ROS_DISTRO={{ ros2_distro }}
Environment=ROBOT_WS={{ workspace_path }}
Environment=HOME=/home/{{ target_user }}
WorkingDirectory=/home/{{ target_user }}
ExecStart=/opt/questix_robot/questix_robot_launcher.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=questix_robot

TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
